# systemd/netflow-dashboard.service
[Unit]
Description=NetFlow Analytics Dashboard
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root
ExecStart=/usr/bin/python3 /root/netflow-dashboard.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
